# images/asteroid/files/system_files/etc/profile.d/asteroid-neofetch.sh
alias neofetch='/usr/bin/fastfetch --color white -c /usr/share/asteroid/configs/fastfetch/config.jsonc'
alias neowofetch='/usr/bin/fastfetch --color white -c /usr/share/asteroid/configs/fastfetch/config.jsonc'
alias fastfetch='/usr/bin/fastfetch --color white -c /usr/share/asteroid/configs/fastfetch/config.jsonc'
